# datasophon-worker/src/main/resources/script/ddh-worker
#!/bin/sh
# chkconfig: 35 95 1
usage="Usage: start.sh (start|stop|restart) <command> "

startStop=$1


case $startStop in
  (start)
    sh /opt/datasophon/ddh-worker-1.0.0/ddh-worker/bin/ddh-worker.sh start worker
    ;;
  (stop)
    sh /opt/datasophon/ddh-worker-1.0.0/ddh-worker/bin/ddh-worker.sh stop worker
      ;;
  (status)
    sh /opt/datasophon/ddh-worker-1.0.0/ddh-worker/bin/ddh-worker.sh status worker
      ;;
  (restart)
    sh /opt/datasophon/ddh-worker-1.0.0/ddh-worker/bin/ddh-worker.sh restart worker
      ;;
  (*)
    echo $usage
    exit 1
    ;;

esac

echo "End $startStop $command."